# contrib/hrs.service
# systemd user service for hrs serve
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp contrib/hrs.service ~/.config/systemd/user/hrs.service
#
# Enable (start on login):
#   systemctl --user enable hrs
#
# Start now:
#   systemctl --user start hrs
#
# Check status:
#   systemctl --user status hrs
#
# View logs:
#   journalctl --user -u hrs -f
#
# Stop and disable:
#   systemctl --user stop hrs
#   systemctl --user disable hrs

[Unit]
Description=hrs time tracking server
After=network.target

[Service]
Type=simple
ExecStart=%h/go/bin/hrs serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
